# docker/supervisor/apache.conf
[program:apache2]
command=apache2-foreground
autostart=true
autorestart=true
killasgroup=true
stopasgroup=true