# docker/supervisor/crond.conf
[program:crond]
command=crond -f
autostart=true
autorestart=true
killasgroup=true
stopasgroup=true